# docker/worker/supervisor.programs.ini
[program:worker]
command=php ./bin/console messenger:consume amqp --time-limit=3600 -v
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0